# ops/forgeclaw-dashboard.service
[Unit]
Description=ForgeClaw Dashboard
After=forgeclaw.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/projects/ForgeClaw/packages/dashboard
ExecStart=/root/.bun/bin/bun run dev
Restart=always
RestartSec=5
Environment=HOME=/root
Environment=PATH=/root/.local/bin:/root/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StandardOutput=append:/root/.forgeclaw/logs/dashboard.log
StandardError=append:/root/.forgeclaw/logs/dashboard.log

[Install]
WantedBy=multi-user.target
